# shinobi-exporter/build/package/debian/shinobi-exporter/lib/systemd/system/shinobi-exporter.service
[Unit]
Description=Shinobi Prometheus exporter
After=network.target

[Service]
EnvironmentFile=-/etc/default/shinobi-exporter
Type=simple
ExecStart=/usr/bin/shinobi-exporter $SHINOBI_EXPORTER_OPTS
Restart=always

[Install]
WantedBy=multi-user.target
